# Installation Files/THE_APP_NAME.conf
[program:APP_NAME]
command=/home/USERNAME/running/APP_NAME  
autostart=true  
autorestart=true  
stderr_logfile=/var/log/APP_NAME.err.log  
stdout_logfile=/var/log/APP_NAME.out.log 